# systemd/epaper-dashboard.service
[Unit]
Description=Waveshare ePaper 10.85 Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/dashboard
ExecStart=/usr/bin/python3 /home/%i/dashboard/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
